# os/linux/openrc/podman-compose/my-alpine/rootfs/etc/profile.d/90my.sh
# 设置全局环境变量
#export JAVA_HOME="/opt/java/openjdk"
#export PATH="/opt/arthas:${JAVA_HOME}/bin:${PATH}"
export PS1='\n\e[1;37m[\e[m\e[1;32m\u\e[m\e[1;33m@\e[m\e[1;35m\H\e[m \e[4m`pwd`\e[m\e[1;37m]\e[m\e[1;36m\e[m\n\$'
alias tailf='tail -f'
alias ls='ls --color=auto'
alias ll='ls -lh'
alias vi=vim
